# _ops/awwan-tour/mkosi.extra/etc/systemd/system/org-awwan-tour.service
[Unit]
Description=org-awwan-tour
After=network-online.target

[Service]
User=awwan
WorkingDirectory=/home/awwan/tour
ExecStart=/usr/local/bin/awwan -address=0.0.0.0:27600 \
  serve /home/awwan/tour
Restart=always
RestartSec=5s

[Install]
WantedBy=default.target
