# _ops/awwan.org/etc/systemd/system/www-awwan.path
[Unit]
Description="Watch www-awwan"

[Path]
PathChanged={{.Val "awwan:www:bin"}}
Unit=systemctl-restart@%p.service

[Install]
WantedBy=multi-user.target
