# _ops/awwan-play/mkosi.extra/etc/systemd/system/org-awwan-play.path
[Unit]
Description="Watch /usr/local/bin/awwan"

[Path]
PathChanged=/usr/local/bin/awwan
Unit=systemctl-restart@%p.service

[Install]
WantedBy=multi-user.target
